# apps_ban_checker.cron
0 * * * * appscontroller /home/appscontroller/app/.venv/bin/python /home/appscontroller/app/apps_ban_checker.py
