# apps_ban_checker.cron
58 */1 * * * appscontroller python /home/appscontroller/appstrafficcontroller/apps_ban_checker.py
